# docker/app/gunicorn.conf
[program:chatbot]
directory=/home/chatbot
command=gunicorn main:app -b localhost:3001
autostart=true
autorestart=true
